# pickle/files/etc/init.d/gpio_bb
#!/bin/sh /etc/rc.common

START=10
STOP=10

start() {        
	mknod /dev/gpio-bb c 180 0
}                 
 
stop() {          
	rm /dev/gpio-bb
}
